# .config/autostart/psensor.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/psensor
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=psensor
Comment=
